# systemd/wifi-speed.service
[Unit]
Description=WiFi speed measurement (wifi-speed)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/opt/wifi-speed
ExecStart=/opt/wifi-speed/.venv/bin/wifi-speed --config /etc/wifi-speed/config.yaml run
Environment=PYTHONUNBUFFERED=1

# Pi Zero W はメモリが少ないため、測定中のスワップを抑える
Nice=10

[Install]
WantedBy=multi-user.target
